# config/chroot_local-includes/usr/lib/systemd/user/tracker-extract-3.service.d/after-gnome-session.conf
# TODO:Trixie: In the tracker-extract version on Trixie, this After=
# dependency exists in the default service file
# (/usr/lib/systemd/user/tracker-extract-3.service), so we can remove
# this file.
[Unit]
After=gnome-session.target
